# systemd/pve-cert-sync.timer
[Unit]
Description=Run PVE certificate sync frequently

[Timer]
OnCalendar=*-*-* 03,09,15,21:00:00
Persistent=true

[Install]
WantedBy=timers.target
