# systemd/pve-cert-sync.service
[Unit]
Description=PVE Certificate Sync Service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/pve-cert-sync
